# assets/runit-userservices.bak/sv/powertop/run
#!/bin/sh

[ -r ./conf ] && . ./conf

exec powertop ${OPTS:-} 2>&1
