# assets/runit-userservices.bak/sv/androtifyd/finish
#!/usr/bin/env bash

echo ":: $(basename $(pwd)) service is finished..."
